# mail-faroe/.config/systemd/user/mbsync.service
[Unit]
Description=Run mbsync
Wants=mbsync.timer

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -a
